# archives/deployarr_v4/includes/rclone/start-media-after-boot.service
[Unit]
Description=start-media-after-boot
After=network.target

[Service]
ExecStart=SMAB-PATH-PLACEHOLDER

[Install]
WantedBy=default.target